# etc/samples/crontab
# 4.4BSD Style crontab 	$MAIL_LIST
# set this crontab
# % crontab crontab
# to verify the corntab
# % crontab -l
#
0 * * * *	$EXECDIR/msend.pl $MLDIR --mladdr=$MAIL_LIST


# 注意：
#
# もし、コマンドラインオプションを使うのなら、まとめ送りでも
# コマンドラインオプションをつけて下さいね。
# もともと"別の"プロセスだからコマンドラインオプションは当然必要です。
# そういう意味では標準のやり方の config.ph に設定はすべて書いておくと
# いう方がよいですね。
# 
# 	--auto
# 0 * * * *	$EXECDIR/msend.pl $MLDIR --mladdr=$MAIL_LIST --auto
#
# 	--autosubject
# 0 * * * *	$EXECDIR/msend.pl $MLDIR --mladdr=$MAIL_LIST --autosubject
#
# 	--autobody
# 0 * * * *	$EXECDIR/msend.pl $MLDIR --mladdr=$MAIL_LIST --autobody
#
